# packages/net-dns/pdnsd/files/systemd/pdnsd.service
[Unit]
Description=pdnsd proxy DNS server

[Service]
Type=forking
ExecStart=/usr/sbin/pdnsd -d -s -t 
PIDFile=/var/run/pdnsd.pid
ExecReload=/usr/sbin/pdnsd-ctl config

[Install]
WantedBy=multi-user.target
