# packages/net-dns/pdnsd/files/systemd/pdnsd.service
# Copyright 2011 Hong Hao <oahong@oahong.me>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=proxy DNS server
After=network.target
Conflicts=bind.service

[Service]
Type=forking
PIDFile=/run/pdnsd.pid
ExecStart=/usr/sbin/pdnsd -d -s -t -p /run/pdnsd.pid
ExecReload=/usr/sbin/pdnsd-ctl config

[Install]
WantedBy=multi-user.target
