# packages/net-dns/pdnsd/files/systemd/pdnsd.service
[Unit]
Description=pdnsd proxy DNS server
After=network.target

[Service]
Type=forking
PIDFile=/var/run/pdnsd.pid
ExecStart=/usr/sbin/pdnsd -d -s -t -p /var/run/pdnsd.pid
ExecReload=/usr/sbin/pdnsd-ctl config

[Install]
WantedBy=multi-user.target
