# packages/net-dns/dnscrypt-proxy/files/systemd/dnscrypt-proxy.service
[Unit]
Description=DNSCrypt proxy
Documentation=man:dnscrypt-proxy(8)
After=network.target iptables.service
Requires=dnscrypt-proxy.socket

[Service]
Type=notify
NonBlocking=true
NotifyAccess=main
User=dnscrypt-proxy
Environment=DNSCRYPT_PROXY_RESOLVER_NAME=cloudns-syd
ExecStart=/usr/sbin/dnscrypt-proxy \
        --resolver-name=${DNSCRYPT_PROXY_RESOLVER_NAME}
Restart=on-abort

[Install]
WantedBy=multi-user.target
Also=dnscrypt-proxy.socket
